# hw04-patroni-ha/ansible/roles/pgbackrest_client/templates/pgbackrest-backup.timer.j2
# Managed by Ansible
[Unit]
Description=Run pgBackRest backup daily (leader only)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
